# glances/rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Glances
# Runs Glances InfluxDB Export
# ==============================================================================
declare -a options

if bashio::config.false 'influxdb.enabled'; then
    exec sleep 86400
fi

options+=(-C /etc/glances.conf)
options+=(--export influxdb)
options+=(--quiet)

options+=(--time "$(bashio::config 'refresh_time')")

if bashio::config.false 'process_info'; then
    options+=(--disable-process)
fi

if bashio::debug; then
    options+=(--debug)
fi

# Interval
sleep "$(bashio::config 'influxdb.interval')"

# Run Glances
exec glances "${options[@]}"

