# glances/rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Glances
# Runs Glances InfluxDB Export
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

declare -a options

if hass.config.false 'influxdb.enabled'; then
    exec sleep 86400
fi

options+=(-C /etc/glances.conf)
options+=(--export influxdb)
options+=(--quiet)

if hass.debug; then
    options+=(--debug)
fi

# Interval
sleep "$(hass.config.get 'influxdb.interval')"

# Run Glances
exec glances "${options[@]}"

