# glances/rootfs/etc/services.d/influxdb/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: Glances
# Runs Glances InfluxDB Export
# ==============================================================================
declare -a options

if bashio::config.false 'influxdb.enabled'; then
    bashio::exit.ok
fi

options+=(-C /etc/glances.conf)
options+=(--quiet)

if bashio::config.equals 'influxdb.version' '2'; then
    options+=(--export influxdb2)
else
    options+=(--export influxdb)
fi

options+=(--time "$(bashio::config 'refresh_time')")

if bashio::config.false 'process_info'; then
    options+=(--disable-process)
fi

if bashio::debug; then
    options+=(--debug)
fi

while true
do
    # Interval
    sleep "$(bashio::config 'influxdb.interval')"

    # Run Glances
    glances "${options[@]}"
done
