# glances/rootfs/etc/services.d/influxdb/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: Glances
# Runs Glances InfluxDB Export
# ==============================================================================
declare -a options
declare -i influxdb_interval influxdb2_interval interval

if bashio::config.false 'influxdb.enabled' && bashio::config.false 'influxdb2.enabled'; then
    bashio::exit.ok
fi

options+=(-C /etc/glances.conf)
if bashio::config.true 'influxdb.enabled'; then
  options+=(--export influxdb)
fi

if bashio::config.true 'influxdb2.enabled'; then
  options+=(--export influxdb2)
fi
options+=(--quiet)

options+=(--time "$(bashio::config 'refresh_time')")

if bashio::config.false 'process_info'; then
    options+=(--disable-process)
fi

if bashio::debug; then
    options+=(--debug)
fi

# Get whichever interval is set or if both influxdb and influxdb2 take the larger
if bashio::config.true 'influxdb.enabled'; then
    influxdb_interval="$(bashio::config 'influxdb.interval')"
fi

if bashio::config.true 'influxdb2.enabled'; then
    influxdb2_interval="$(bashio::config 'influxdb2.interval')"
fi

if (( ${influxdb_interval:-0} > ${influxdb2_interval:-0} )); then
  interval=${influxdb_interval}
else
  interval=${influxdb2_interval}
fi

while true
do
    # Interval
    sleep "${interval}"

    # Run Glances
    glances "${options[@]}"
done
