# glances/rootfs/etc/services.d/glances/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: Glances
# Runs Glances
# ==============================================================================
declare -a options

bashio::log.info 'Starting Glances...'

options+=(-C /etc/glances.conf)
options+=(-w)
options+=(--port 61209)
options+=(--bind 127.0.0.1)
options+=(--time "$(bashio::config 'refresh_time')")

if bashio::config.false 'process_info'; then
    options+=(--disable-process)
fi

if bashio::debug; then
    options+=(--debug)
fi

# Run Glances
exec glances "${options[@]}"
